# root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

foreground { s6-echo "Starting nginx..." }

foreground { s6-chmod 707 /dev/stderr }

foreground { s6-chmod 707 /dev/stdout }

s6-setuidgid nobody nginx
